# deployment/dev/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
loglevel=info

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================
# MCP Server with Hot Reload (Local Development)
# ============================================
[program:mcp_server]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8081 --reload --log-level info
directory=/Users/xenodennis/Documents/Fun/isA_MCP
autostart=true
autorestart=true
startsecs=5
startretries=3
# Output to stdout/stderr for easy debugging
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/Users/xenodennis/Documents/Fun/isA_MCP",SERVICE_NAME="mcp_server"
priority=10
